# dist/macos/com.chiwanpark.hsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chiwanpark.hsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/${USER}/.local/bin/hsync</string>
        <string>client</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/hsync.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/hsync.err</string>
</dict>
</plist>
